# contrib/systemd/pping.service
[Unit]
Description=pping passive RTT monitor
Documentation=https://github.com/nshopik/pping
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/default/pping
ExecStart=/usr/local/bin/pping-supervise.sh
ExecReload=/bin/kill -HUP $MAINPID
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
