# contrib/clickhouse/pping2-load.service
[Unit]
Description=Batch-load pping2 output into ClickHouse
After=pping2.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pping2-load.sh
