# contrib/clickhouse/pping2-load.cron
# Batch-load pping2 output into ClickHouse every minute.
* * * * * root /usr/local/bin/pping2-load.sh
