# contrib/clickhouse/pping-load.cron
# Batch-load pping output into ClickHouse every minute.
* * * * * root /usr/local/bin/pping-load.sh
